# dot_config/systemd/user/kanata-keychron-k2-he.service
[Unit]
Description=Kanata keyboard remapper (Keychron K2 HE)
Documentation=https://github.com/jtroo/kanata

[Service]
Environment=PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/bin
Environment=DISPLAY=:0
Environment=HOME=/home/stephen
Type=simple
ExecStart=/home/stephen/.local/bin/kanata --cfg /home/stephen/.config/kanata/keychron-k2-he.kbd
Restart=never

[Install]
WantedBy=default.target
